# packages/tagihan/files/etc/init.d/tagihan
#!/bin/sh /etc/rc.common
START=99
STOP=10
CRON_TAG="tagihan-cron"

start() {
    local interval=$(uci -q get tagihan.settings.interval)
    [ -z "$interval" ] && interval=6
    crontab -l 2>/dev/null | grep -v "$CRON_TAG" > /tmp/cron_tmp
    echo "0 */${interval} * * * /usr/bin/tagihan-check all # ${CRON_TAG}" >> /tmp/cron_tmp
    crontab /tmp/cron_tmp
    rm -f /tmp/cron_tmp
    /usr/bin/tagihan-check all > /tmp/tagihan-last.log 2>&1 &
    echo "Tagihan Bot dimulai. Interval: setiap ${interval} jam."
}

stop() {
    crontab -l 2>/dev/null | grep -v "$CRON_TAG" | crontab -
    echo "Tagihan Bot dihentikan."
}

restart() { stop; sleep 1; start; }
